# scripts/launchd/com.smithers.fathom-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd plist for the Smithers Fathom Sync.

  Hits POST http://localhost:3000/api/jobs/fathom-sync on the configured
  interval. Fires whether or not `pnpm dev` is running — when dev is down
  the curl will fail silently (logged to /tmp/smithers-fathom-sync.log).

  To install:
    1. cp this file to ~/Library/LaunchAgents/com.smithers.fathom-sync.plist
       (drop the .example suffix)
    2. Edit StartInterval below to match your /settings interval (in seconds).
    3. launchctl load ~/Library/LaunchAgents/com.smithers.fathom-sync.plist

  To uninstall:
    launchctl unload ~/Library/LaunchAgents/com.smithers.fathom-sync.plist
    rm ~/Library/LaunchAgents/com.smithers.fathom-sync.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smithers.fathom-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/curl</string>
        <string>-sS</string>
        <string>-X</string>
        <string>POST</string>
        <string>--max-time</string>
        <string>60</string>
        <string>http://localhost:3000/api/jobs/fathom-sync</string>
    </array>

    <!-- Seconds between fires. 3600 = 1 hour. Match your /settings interval. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/smithers-fathom-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/smithers-fathom-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
